# optimizePNG_supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:optimizer]
command=/bin/sh -c "/opt/scripts/optimizePNG.sh"
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
