# optimizePNG_supervisord.conf
[supervisord]
nodaemon=true

[program:optimizer]
command=/bin/sh -c "/opt/scripts/optimizePNG.sh"
autorestart=true
