# novnc/conf.d/desktop.conf
[program:xfce4]
command=startxfce4
autorestart=true