# novnc/conf.d/websockify.conf
[program:websockify]
command=websockify --web /usr/share/novnc 8088 localhost:5900
autorestart=true
